# phase-3-arch-postinstall/metis-ddns/metis-ddns.timer
[Unit]
Description=metis dynamic DNS update — periodic re-check

[Timer]
# 2 min after boot, then every 10 min. Cheap when nothing changed (one
# ipify probe per stack, no Azure traffic), so the 10-min cadence is fine.
# Persistent=true catches up if the laptop was sleeping past a tick.
OnBootSec=2min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
