# phase-3-arch-postinstall/system-files/udev/99-tablet-mode.rules
# /etc/udev/rules.d/99-tablet-mode.rules
#
# Inspiron 7786 tablet-mode signaling:
#
# The 7786 has NO dedicated SW_TABLET_MODE hinge sensor. Tablet-mode events
# come from ACPI `intel-vbtn` (codes 0xCC enter / 0xCD exit), which the
# kernel's intel_vbtn driver synthesizes onto a virtual input device
# named exactly "Intel HID events" under /devices/platform/INT33D5:00/.
#
# That device emits EV_SW / SW_TABLET_MODE events, and udev relays them as
# `change` events with ENV{SW_TABLET_MODE}={0,1} on the parent input
# device. Catching those gives us a clean userspace edge.
#
# Why we write to a state file instead of running hyprctl directly:
#   - udev runs as root with no Wayland/Hyprland env. Calling hyprctl from
#     here would need a wrapper that finds WAYLAND_DISPLAY +
#     HYPRLAND_INSTANCE_SIGNATURE and `runuser` into tom — fragile across
#     login sessions, and breaks at the greeter.
#   - Instead: we write the new SW_TABLET_MODE bit to /run/tablet-mode/state.
#     A user-level systemd path unit (~/.config/systemd/user/tablet-mode-watch.path)
#     watches that file and triggers tablet-mode.service, which inherits
#     the user session env and runs hyprctl with full context.
#
# The mkdir -p cost of doing this on every event is negligible (<1ms),
# and the unconditional `printf` keeps the rule branch-free.
#
# Verify on real hardware:
#   sudo udevadm monitor --environment --subsystem-match=input
#   (then fold the screen — should see ACTION=change with SW_TABLET_MODE=1)

ACTION=="change", SUBSYSTEM=="input", ATTRS{name}=="Intel HID events", \
    RUN+="/bin/sh -c 'mkdir -p /run/tablet-mode && printf %%s %E{SW_TABLET_MODE} > /run/tablet-mode/state && chmod 0644 /run/tablet-mode/state'"
