# dotfiles/dot_config/systemd/user/tablet-mode.service
[Unit]
Description=Apply current tablet/laptop mode (kbd + touchpad + OSK)
# Only meaningful inside a Wayland session — at greeter time there's no
# Hyprland to talk to, so skip cleanly.
ConditionEnvironment=WAYLAND_DISPLAY
# Run sequentially with itself: if the lid is folded twice in <1s the
# second invocation queues up rather than racing the first.
RefuseManualStart=no

[Service]
Type=oneshot
ExecStart=%h/.local/bin/tablet-mode-toggle --detect
# tablet-mode-toggle reads /run/tablet-mode/state — written by the udev
# rule that triggered the path unit. --detect handles the case where
# the file doesn't exist yet (cold-boot fallback to laptop mode).

[Install]
WantedBy=default.target
