# dotfiles/dot_config/systemd/user/wallpaper-rotate.service
[Unit]
Description=Rotate wallpaper + regenerate matugen palette
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
Type=oneshot
ExecStart=%h/.local/bin/wallpaper-rotate --auto

[Install]
WantedBy=default.target
